# upgrade_mailerpy/cfg/cron.daily/upgrade_mailer
#!/bin/sh

exec /usr/local/bin/upgrade_mailer.py &> /dev/null
